# user/eating.target
[Unit]
Description=Eating...
Conflicts=thinking.target hungry.target
Wants=eating.timer eating-food.service
