# user/hungry.target
[Unit]
Description=Hungry...
Conflicts=eating.target thinking.target
Wants=hungry.timer hunger.service
