# user/select-seat.path
[Unit]
Description=Watch Dining Room

[Path]
PathModified=/home/share/dining-room

[Install]
WantedBy=default.target
