# user/thinking.target
[Unit]
Description=Thinking...
Conflicts=eating.target hungry.target
Wants=contemplation.service thinking.timer
